# program/wd.cron
#Watchdog tasks 

* * * * * /etc/init.d/schedule_daemon_wd
* * * * * /etc/init.d/schedule_webserver_wd